# schedule.cron
# Run the affiliate automation every day at 06:00
# Edit with: crontab -e
0 6 * * * cd /workspace && /usr/bin/python3 run.py >> /workspace/run.log 2>&1
